# hot_tub/crontab
0 0 * * * /sbin/reboot 
@reboot cd /path/to/cabinpi/hot_tub && ./on_startup
